# deploy/systemd/reeve-daemon.service.template
[Unit]
Description=Reeve Pulse Daemon
After=network.target
Wants=network-online.target

[Service]
Type=simple
User={{USER}}
Group={{USER}}
WorkingDirectory={{REEVE_BOT_PATH}}

# Use uv to run the daemon
ExecStart={{UV_PATH}} run python -m reeve.pulse

# Environment
EnvironmentFile={{REEVE_BOT_PATH}}/.env
# PATH must include node (for hapi) and claude CLI locations
# NOTE: Adjust the nvm path for your node version
Environment="PATH=/home/{{USER}}/.local/bin:/home/{{USER}}/.nvm/versions/node/v18.20.7/bin:/usr/local/bin:/usr/bin:/bin"

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Logging (stdout/stderr -> journald)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=reeve-daemon

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths={{REEVE_HOME}} {{REEVE_DESK_PATH}} /home/{{USER}}/.cache/uv /home/{{USER}}/.hapi /home/{{USER}}/.claude {{EXTRA_RW_PATHS}}

# Resource limits (Claude sessions need significant memory)
MemoryMax=2G
CPUQuota=100%

[Install]
WantedBy=multi-user.target
